# debian/69-nvidia-persistenced.rules
ACTION=="add" DEVPATH=="/module/bbswitch" SUBSYSTEM=="module" RUN+="/usr/bin/start-nvidia-persistenced"
